# deploy/devflow-worker.service
# Systemd unit for the devflow worker.
#
# Install (once, on the server):
#   sudo cp deploy/devflow-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now devflow-worker
#
# Replace REPLACE_ME with the deploy account username. Ensure `pnpm` is on the
# login shell PATH (e.g. `corepack enable`) — ExecStart uses bash -l for that.

[Unit]
Description=devflow worker (claims jobs from the queue)
After=network.target

[Service]
Type=simple
User=REPLACE_ME
WorkingDirectory=/home/REPLACE_ME/devflow
EnvironmentFile=/home/REPLACE_ME/devflow/.env
ExecStart=/bin/bash -lc 'pnpm worker'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
